# .onreboot.crontab
# "On reboot - run ensure Kratos is up and running."
@reboot cd /workspace/tjmaynes/kratos && make start