# cron.d/onreboot.crontab
# "On reboot - ensure GECK is up and running after reboot."
@reboot cd /workspace/tjmaynes/geck && make start